# notif_app/myjob.cron
59 12 * * 1-5 /bin/bash /Users/diego/Projects/get_weather/morning_call.sh > /dev/null 2>&1
50 17 * * 1-5 /bin/bash /Users/diego/Projects/get_weather/morning_call.sh > /dev/null 2>&1

